# config/systemd/user/bluetooth-autoconnect.service
[Unit]
Description=Auto-conectar Alexa via Bluetooth
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/bin/conectar_bluetooth.sh
# trata códigos não-críticos como sucesso (ex.: 2 = MAC ausente, 3 = comando não encontrado)
SuccessExitStatus=0 2 3
