# config/systemd/user/hyprpaper-change.service
[Unit]
Description=Troca automática de wallpapers no Hyprland

[Service]
ExecStart=%h/.config/scripts/hyprpaper_change.sh